# deploy/mnemosyne-daemon.service
# Mnemosyne auto-save daemon — sweeps agent transcripts into the palace.
#
#   sudo cp deploy/mnemosyne-daemon.service /etc/systemd/system/
#   systemctl --user daemon-reload && systemctl --user enable --now mnemosyne-daemon
#
# Runs as a user service so it can read ~/.claude/projects.

[Unit]
Description=Mnemosyne transcript auto-save daemon

[Service]
Type=simple
ExecStart=%h/.cargo/bin/mnemosyne daemon run --watch %h/.claude/projects --interval 300
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
